# netwatchm-web.service
[Unit]
Description=NetWatchM Web Dashboard
Documentation=file:///var/lib/netwatchm/report.html
After=network.target netwatchm.service
Wants=netwatchm.service
OnFailure=netwatchm-notify@%n.service

[Service]
Type=simple
WorkingDirectory=/var/lib/netwatchm
ExecStart=/usr/bin/python3 -m http.server 8765
Restart=always
RestartSec=5
StartLimitIntervalSec=120
StartLimitBurst=5
User=root

[Install]
WantedBy=multi-user.target
