# netwatchm-web.service
[Unit]
Description=NetWatchM Web Dashboard
Documentation=file:///var/lib/netwatchm/report.html
After=network.target netwatchm.service
Wants=netwatchm.service
OnFailure=netwatchm-notify@%n.service
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/var/lib/netwatchm
Environment=NETWATCHM_CMD=/home/jbaez120/.local/bin/netwatchm
Environment=NETWATCHM_GEOIP_DB=/var/lib/netwatchm/GeoLite2-City.mmdb
Environment=NETWATCHM_FLOW_DB=/var/lib/netwatchm/flows.db
ExecStart=/usr/local/bin/netwatchm-server
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
